# profile.d/init.sh
export LANG=en_US.UTF-8

PS1="[\u@\h \W\[\e[31m\] \$(git-branch-prompt)\[\e[0m\]] $ "
